# toddleglow.service
[Unit]
Description=ToddleGlow

[Service]
User=pi
ExecStart=/home/pi/toddleglow/toddleglow.py

[Install]
WantedBy=multi-user.target
